# dist/proton-fuse.service
[Unit]
Description=Proton per-user FUSE filesystem
After=default.target

[Service]
Type=notify
ExecStart=/usr/local/bin/proton-fuse
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
