# dist/protonfs.service
[Unit]
Description=Proton per-user FUSE filesystem
Documentation=https://github.com/major0/proton-utils
After=protonfs-redirector.service

[Service]
Type=notify
ExecStart=/usr/local/sbin/proton-fuse
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
